# docker/data/supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=/usr/bin/tor --defaults-torrc /usr/share/tor/tor-service-defaults-torrc
stderr_logfile = /var/log/supervisor/tor-stderr.log
stdout_logfile = /var/log/supervisor/tor-stdout.log

[program:globaleaks]
command=/usr/bin/python3 /usr/bin/globaleaks --ip=0.0.0.0 --user=globaleaks --group=globaleaks --working-path=/var/globaleaks/ -n
stderr_logfile = /var/log/supervisor/globaleaks-stderr.log
stdout_logfile = /var/log/supervisor/globaleaks-stdout.log